# db/hockey-db-update.timer
[Unit]
Description=Run the Hockey DB update service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target